# registry/templates/docker-registry.conf.erb
limit nofile 10000 10000

kill timeout 180

start on runlevel [2345]
stop on runlevel [06]

setuid <%= @user %>
setgid <%= @group %>

respawn

env DOCKER_REGISTRY_CONFIG=<%= @path %>/config.yml
script
<%= @venv_path %>/bin/gunicorn \
    -k gevent --max-requests <%= @gunicorn_max_requests %> \
    --graceful-timeout 3600 -t 3600 -b <%= @ipbind_port %> \
    -w <%= @gunicorn_workers %> docker_registry.wsgi:application
end script
post-stop exec sleep 5
